# config/schedules/closed-loop-production-daily.cron.example
CRON_TZ=Europe/Oslo
30 5 * * * /home/hermoine/agent-research-book-pipeline/scripts/run_production_daily_cron.sh
